# config/tty.service
[Service]
ExecStart=
ExecStart=/sbin/mingetty --autologin gns3 --noclear %I
